# docker/cron/crontab
0 12 * * * php /var/www/html/yii recurring-transaction/generate >> /var/www/html/runtime/logs/cron.log 2>&1
